# roles/2-common/templates/iiab-startup.service
[Unit]
Description=Execute startup script  
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libecec/iiab-startup.sh

[Install]
WantedBy=multi-user.target
