# roles/network/templates/hostapd/70-persistent-net.rules
SUBSYSTEM=="ieee80211", ACTION=="add|change", KERNEL=="phy0", \
  RUN+="/sbin/iw phy phy0 interface add ap0 type __ap", \
  RUN+="/sbin/ip link set ap0 address {{ ap0_mac_addr }}"
