# roles/network/templates/hostapd/iiab-wifi-test.service.j2
[Unit]
Description=IIAB find channel freq for ssid
Requires=network-pre.target
BindsTo=sys-subsystem-net-devices-{{ discovered_wireless_iface }}.device
After=sys-subsystem-net-devices-{{ discovered_wireless_iface }}.device
After=network-pre.target

After=iiab-clone-wifi.service
After=wpa_supplicant.service
Wants=wpa_supplicant.service
After=NetworkManager.service
After=netplan-wpa-{{ discovered_wireless_iface }}.service
Before=hostapd.service
Before=dhcpcd.service
Before=wpa_supplicant@{{ discovered_wireless_iface }}.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/iiab-test-wifi

[Install]
WantedBy=multi-user.target
