# roles/bluetooth/templates/bt-pan.service.j2
[Unit]
Description=Bluetooth NEP PAN
After=br0.network

[Service]
ExecStart=/usr/bin/bt-network -s nap br0
Type=simple

[Install]
WantedBy=multi-user.target

