# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
User={{ calibre_web_user }}
ExecStart=/usr/bin/python " {{ calibre_web_exec_path }} "
WorkingDirectory={{ calibre_web_path }}

[Install]
WantedBy=multi-user.target
