# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
User={{ calibreweb_user }}
#EnvironmentFile={{ calibreweb_home }}/cps.conf
#WorkingDirectory={{ calibreweb_path }}
ExecStart=/usr/bin/python {{ calibreweb_exec_path }} -p {{ calibreweb_home }}/{{ calibreweb_settings_database }}


[Install]
WantedBy=multi-user.target
