# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
Environment=CALIBRE_DBPATH={{ calibreweb_home }}
User={{ calibreweb_user }}
ExecStart={{ calibreweb_exec_path }} -p {{calibreweb_home}}{{ calibreweb_database }}
WorkingDirectory={{ calibreweb_path }}

[Install]
WantedBy=multi-user.target
