# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
User={{ calibre-web_user }}
ExecStart=/usr/bin/python " {{ calibre-web_exec_path }} "
WorkingDirectory={{ calibre-web_path }}

[Install]
WantedBy=multi-user.target
