# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
User={{ calibreweb_user }}
ExecStart={{ calibreweb_exec_path }} -p {{ calibreweb_database }}
WorkingDirectory={{ calibreweb_path }}

[Install]
WantedBy=multi-user.target
