# roles/calibre-web/templates/calibre-web.service.j2
Description=Calibre-Web

[Service]
Type=simple
User={{ calibreweb_user }}
ExecStart={{ calibreweb_exec_path }} -p {{ calibreweb_config }}/{{ calibreweb_settings_database }}

[Install]
WantedBy=multi-user.target
