# roles/sugar-stats/templates/statistics-consolidation/stats-consolidation.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 3 * * * root stats_consolidation_run --config_file /etc/stats-consolidation.conf

