# roles/1-prep/templates/92-rtc-i2c.rules
# /etc/udev/rules.d/92-rtc-i2c.rules
#
ACTION=="add", SUBSYSTEM=="rtc", ATTRS{hctosys}=="0", RUN+="/usr/sbin/hwclock -s --utc"
