# roles/usb_lib/templates/usbmount@.service.j2
[Unit]
BindTo=%i.device
After=%i.device
After=rc-local.service

[Service]
Type=oneshot
TimeoutStartSec=0
Environment=DEVNAME=%I
ExecStart=/usr/share/usbmount/usbmount add
RemainAfterExit=yes

