# roles/moodle-1.9/moodle/templates/moodle-xs.service.j2
[Unit]
Description=Moodle Course Management OLPC edition
After={{ apache_service }}.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/moodle-xs-init start

[Install]
WantedBy=multi-user.target
