# roles/moodle-1.9/moodle/templates/moodle-xs.service.in
[Unit]
Description=Moodle Course Management OLPC edition
After={{ apache_service }}.service

[Service]
Type=oneshot
{% if is_debian %}
ExecStart=/etc/init.d/moodle-xs-init start
{% else $}
ExecStart=/usr/libexec/moodle-xs-init start
{% endif %}

[Install]
WantedBy=multi-user.target
