# roles/lokole/templates/lokole_celery_beat.conf
[program:lokole_celery_beat]
command={{ lokole_venv }}/bin/celery --app=opwen_email_client.webapp.tasks beat --pidfile={{ lokole_run_dir }}/lokole_celery_beat.pid --loglevel=error
autostart=true
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile={{ lokole_log_dir }}/lokole_celery_beat.stderr.log
stdout_logfile={{ lokole_log_dir }}/lokole_celery_beat.stdout.log
user={{ lokole_user }}
environment=OPWEN_SETTINGS={{ lokole_settings }}
