# roles/lokole/templates/restarter.service.j2
[Unit]
Description=Lokole webapp restarter
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash {{ lokole_run_directory }}/webapp-restart.sh

[Install]
WantedBy=multi-user.target
