# roles/magrit/templates/magrit.service.j2
[Unit]
Description=Magrit
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/magrit/bin/magrit -p {{ magrit_port }}

[Install]
WantedBy=multi-user.target

