# roles/rapidpro/templates/rapidpro-mailroom.service.j2
[Unit]
Description=RapidPro Mailroom Service
After=network.target postgresql.service {{ valkey_service_name }}

[Service]
User={{ rapidpro_user }}
Group={{ rapidpro_user }}
WorkingDirectory={{ rapidpro_dir }}

# Use Valkey connection scheme as preferred by the application code
Environment='MAILROOM_DB=postgres://{{ db_user }}:{{ db_pass }}@localhost/{{ db_name }}?sslmode=disable'
Environment='MAILROOM_VALKEY=valkey://localhost:6379/0'
Environment='MAILROOM_SESSION_STORAGE=db'
Environment='MAILROOM_SPOOL_DIR={{ rapidpro_dir }}/spool/mailroom'
Environment='MAILROOM_DYNAMO_ENDPOINT=http://localhost:1'
Environment='MAILROOM_S3_ENDPOINT=http://localhost:1'
Environment='MAILROOM_ELASTIC='

ExecStart=/usr/local/bin/rapidpro-mailroom
Restart=on-failure

[Install]
WantedBy=multi-user.target
