# roles/rapidpro/templates/rapidpro-mailroom.service.j2
[Unit]
Description=RapidPro Mailroom Service
After=network.target postgresql.service {{ valkey_service_name }}

[Service]
User={{ iiab_admin_user }}
Group={{ iiab_admin_user }}
WorkingDirectory={{ rapidpro_dir }}

# Use Valkey connection scheme as preferred by the application code
Environment='MAILROOM_DB=postgres://{{ rapidpro_db_user }}:{{ rapidpro_db_pass }}@localhost/{{ rapidpro_db_name }}?sslmode=disable'
Environment='MAILROOM_VALKEY=valkey://localhost:6379/0'
Environment='MAILROOM_SESSION_STORAGE=db'
Environment='MAILROOM_SPOOL_DIR={{ rapidpro_dir }}/spool/mailroom'
Environment='AWS_ACCESS_KEY_ID=dummy'
Environment='AWS_SECRET_ACCESS_KEY=dummy'
Environment='AWS_REGION=us-east-1'
Environment='MAILROOM_DYNAMO_ENDPOINT=http://localhost:6000'
Environment='MAILROOM_S3_ENDPOINT=http://localhost:1'
Environment='MAILROOM_ELASTIC='

ExecStart=/usr/local/bin/mailroom
Restart=on-failure

[Install]
WantedBy=multi-user.target
